# templates/systemd/set-web.service
[Unit]
Description=SET Web Dashboard
Documentation=https://github.com/tatargabor/set-core
After=default.target

[Service]
ExecStart=%h/.local/bin/set-orch-core serve --port 7400
Restart=always
RestartSec=5
Environment=PYTHONPATH=__SET_TOOLS_ROOT__/lib:__SET_TOOLS_ROOT__/modules/web
Environment=PATH=%h/.local/bin:/home/linuxbrew/.linuxbrew/bin:/usr/local/bin:/usr/bin:/bin

[Install]
WantedBy=default.target
